# config/systemd/lobby-auto-update.service
[Unit]
Description=Automatic Lobby System Updates
# Removed network-online.target dependency to avoid boot delays
# This service should only run via timer, not during boot

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/lobby-auto-update.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=3600

[Install]
# Removed WantedBy to prevent auto-start during boot
# Service only runs via timer scheduling
EOF < /dev/null