# config/systemd/plymouth-quit-wait.service
[Unit]
Description=Hold until boot process finishes up (Kiosk Version)
After=lobby-compositor.service
Wants=lobby-compositor.service

[Service]
Type=oneshot
# Quit Plymouth immediately after compositor starts
ExecStart=/bin/bash -c "plymouth quit"
RemainAfterExit=yes
TimeoutStartSec=10
[Install]
WantedBy=graphical.target
